# systemd/gerenuk-cleaner.timer
[Unit]
Description=Gerenuk cleaner timer

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
